# configs/launchagents/scalp_batch7/com.liangh.scalp.batch7.mm.ethusdt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.liangh.scalp.batch7.mm.ethusdt</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/wl/projects/LiangH/scripts/run_scalp_batch7_market_maker.sh</string>
    <string>/Users/wl/projects/LiangH/configs/market_maker/hft_inventory_mm_batch7_ethusdt.json</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/wl/projects/LiangH</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>SCALP_BATCH7_MM_DURATION_SECONDS</key>
    <string>86400</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/wl/projects/LiangH/output/fleet/scalp_suite_batch7_24bot/logs/hft_inventory_mm_batch7_ethusdt.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/wl/projects/LiangH/output/fleet/scalp_suite_batch7_24bot/logs/hft_inventory_mm_batch7_ethusdt.launchd.err.log</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
